# directslave/etc/rc.d/freebsd.directslave
#!/bin/sh
#
# PROVIDE: directslave
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name="directslave"
rcvar=directslave_enable

load_rc_config $name

: ${directslave_enable="NO"}
: ${directslave_flags="--run"}
: ${directslave_pidfile="/usr/local/directslave/run/directslave.pid"}

stop_postcmd="directslave_post"

command="/usr/local/directslave/bin/directslave"

pidfile=${directslave_pidfile}

run_rc_command "$1"
